# halium10-post-install.service
[Unit]
Description=Halium 10 Post Install.

[Service]
Type=simple
ExecStart=/usr/bin/halium10-post-install

[Install]
WantedBy=default.target
